# src/packer/postgres/conf/supervisor/postgres.conf
[program:postgres]
command=/usr/lib/postgresql/10/bin/pg_ctl -w -o "-h %(ENV_HOSTNAME)s" -D /var/pgsql_data -l /var/log/postgres/log start
user=postgres-admin
